# includes/zulip/files/supervisor/conf.d/cron.conf
[program:cron]
command		= /usr/sbin/cron -f
stdout_logfile	= /var/log/supervisor/%(program_name)s.log
stderr_logfile	= /var/log/supervisor/%(program_name)s.log
autorestart	= true
