# includes/zulip-puppet/files/supervisor/conf.d/nginx.conf
[program:cron]
command = /usr/sbin/nginx
stdout_events_enabled=true
stderr_events_enabled=true
autorestart = true
