# includes/zulip-puppet/files/supervisor/conf.d/rabbitmq-server.conf
# Remove when zulip has RABBIT_HOST support
[program:cron]
command = /usr/sbin/rabbitmq-server
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true
