# includes/zulip-puppet/files/supervisor/conf.d/zulip_postsetup.conf
[program:zulip-postsetup-create_user]
command = bash -c "sleep 8;/opt/setupZulipUser.sh"
user=zulip
stdout_events_enabled=true
stderr_events_enabled=true
autorestart = false
