# includes/zulip-puppet/files/supervisor/conf.d/zulip_postsetup.conf
[program:zulip-postsetup]
command = /usr/local/bin/zulipPostSetup.sh
stdout_events_enabled=true
stderr_events_enabled=true
autorestart = true
