# includes/zulip-puppet/files/supervisor/conf.d/zulip_postsetup.conf
[program:zulip-postsetup]
command = bash -c "/usr/local/bin/zulipCreateUser.sh"
stdout_events_enabled=true
stderr_events_enabled=true
autorestart = false
