# includes/supervisor/conf.d/zulip_postsetup.conf
[program:zulip-postsetup-create_user]
command = bash -c "sleep 2;/opt/createZulipAdmin.sh"
user=root
stdout_events_enabled=true
stderr_events_enabled=true
autorestart = false
exitcodes = 0,1,127
